# vps/spheris-stream.service
[Unit]
Description=Spheris Stream Auth Server
After=network.target nginx.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/stream/app
Environment="SPHERIS_ADMIN_KEY=spheris-admin-dev"
ExecStart=/var/www/stream/app/venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
